# systemd/mdk-listener-analytics.service
[Unit]
Description=Sample aggregate MDK Space Radio listener analytics
After=network-online.target icecast2.service
Wants=network-online.target

[Service]
Type=oneshot
User=diego
Group=diego
WorkingDirectory=/opt/swr-radio
RuntimeDirectory=mdk-listener-analytics
RuntimeDirectoryMode=0750
UMask=0027
ExecStart=/usr/bin/flock --nonblock /run/mdk-listener-analytics/sampler.lock /opt/swr-radio/bin/sample-listener-analytics.py --database /opt/swr-radio/data/analytics.sqlite3 --status-url http://127.0.0.1:8000/status-json.xsl --mount /strange-waves.mp3 --timeout 8
ExecStartPost=/opt/swr-radio/bin/generate-listener-analytics.py --database /opt/swr-radio/data/analytics.sqlite3 --output /home/mdk/web/radio.mdkband.com/public_html/stats/listener-analytics.json --mount /strange-waves.mp3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/swr-radio/data /home/mdk/web/radio.mdkband.com/public_html/stats /run/mdk-listener-analytics
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
