# systemd/mdk-play-history-consumer.timer
[Unit]
Description=Consume MDK Space Radio play-history spool every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Persistent=true
Unit=mdk-play-history-consumer.service

[Install]
WantedBy=timers.target
