# systemd/mdk-listener-analytics.timer
[Unit]
Description=Sample MDK Space Radio aggregate listeners every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=true
Unit=mdk-listener-analytics.service

[Install]
WantedBy=timers.target
